# systemd/shotclock-kiosk.service
[Unit]
Description=Shotclock Pi Kiosk Service
Documentation=https://docs.shotclock.local
After=network.target shotclock-agent.service
Requires=shotclock-agent.service

[Service]
Type=simple
# Chromium must attach to the active graphical login session. The launcher runs
# as root only long enough to find that session, then drops to KIOSK_USER or the
# active desktop user before starting Chromium.
User=root
Group=root
WorkingDirectory=/opt/shotclock/current
EnvironmentFile=-/opt/shotclock/shared/.env
ExecStart=/bin/bash /opt/shotclock/current/scripts/launch-kiosk.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shotclock-kiosk

# The kiosk needs access to display
Environment=DISPLAY=:0

# Hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/opt/shotclock /home/shotclock/.shotclock
PrivateTmp=true

[Install]
WantedBy=multi-user.target
