# systemd/shotclock-agent.service
[Unit]
Description=Shotclock Pi Agent Service
Documentation=https://docs.shotclock.local
After=network.target

[Service]
Type=simple
# The agent manages WiFi AP setup, dnsmasq/hostapd, and iptables. Keep HOME
# pointed at the shotclock user's state directory while allowing root-only setup.
User=root
Group=root
WorkingDirectory=/opt/shotclock/current
Environment=NODE_ENV=production
Environment=HOME=/home/shotclock
EnvironmentFile=/opt/shotclock/shared/.env
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/node /opt/shotclock/current/apps/pi-agent/dist/index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shotclock-agent

# Hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/opt/shotclock /home/shotclock/.shotclock
PrivateTmp=true

[Install]
WantedBy=multi-user.target
